# service/apprun/run
#!/bin/bash -e

# Make sure we're in the right place
cd /app


export HOME=/app

# Run the actual app here
	exec 2>&1 \
	     /app/bin/steemyo 2>&1
