# service/yo/run
#!/bin/bash -e

# Make sure we're in the right place
cd "${APP_ROOT}"


export HOME=/app

# Run the actual app here
    exec 2>&1 \
         pipenv run python -m yo.cli 2>&1
